# dns-update-infomaniak/dns-update-infomaniak.timer
[Timer]
OnCalendar=Minutely
Unit=dns-update-infomaniak.service

[Install]
WantedBy=timers.target
