# dns-update-infomaniak/dns-update-infomaniak.service
[Service]
Type=simple
ExecStart=/usr/local/bin/dns-update-infomaniak
